# figma-agent.socket
[Unit]
Description=Figma Agent for Linux

[Socket]
ListenStream=44950

[Install]
WantedBy=sockets.target
